# scripts/services/cdu-local-ui.service
[Unit]
Description=L2A CDU Local UI (X server + PySide6 kiosk)
After=multi-user.target redis.service systemd-user-sessions.service
Wants=redis.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=gadgetini
PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/gadgetini/L2A_CDU_system

Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=QT_QPA_PLATFORM=xcb
Environment=QT_SCALE_FACTOR=1

ExecStart=/usr/bin/xinit /home/gadgetini/L2A_CDU_system/scripts/cdu_session.sh -- :0 vt1 -nocursor -nolisten tcp
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
