# scripts/services/cdu-mcg.service
[Unit]
Description=L2A CDU MCG (Modbus Control Gateway)
After=redis.service
Wants=redis.service

[Service]
Type=simple
User=gadgetini
WorkingDirectory=/home/gadgetini/L2A_CDU_system
ExecStart=/home/gadgetini/venv/bin/python -m src.mcg.main
# Allow access to all common Modbus serial device groups:
#   dialout - /dev/ttyACM* (USB CDC), /dev/ttyAMA0 (= /dev/serial0)
#   plugdev - /dev/ttyUSB* (USB-UART adapters like FTDI on some RPi images)
SupplementaryGroups=dialout plugdev
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
