# scripts/services/cdu-web-backend.service
[Unit]
Description=L2A CDU Web UI backend (FastAPI + SvelteKit static bundle)
After=redis.service
Wants=redis.service

[Service]
Type=simple
User=gadgetini
WorkingDirectory=/home/gadgetini/L2A_CDU_system
ExecStart=/home/gadgetini/venv/bin/uvicorn src.web_ui.backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
